# packaging/systemd/oc-rsyncd.service
# packaging/systemd/oc-rsyncd.service
[Unit]
Description=oc-rsync daemon
Documentation=man:oc-rsyncd(8) man:oc-rsyncd.conf(5) man:oc-rsync(1)
Wants=network-online.target
After=network-online.target

[Service]
Type=notify
NotifyAccess=main
User=ocrsync
Group=ocrsync
RuntimeDirectory=oc-rsyncd
LogsDirectory=oc-rsyncd
StateDirectory=oc-rsyncd
ConfigurationDirectory=oc-rsyncd
ExecStart=/usr/bin/oc-rsync --daemon --no-detach --config=/etc/oc-rsyncd.conf
# To run the daemon via the dedicated oc-rsyncd binary instead of the main oc-rsync CLI,
# create a drop-in with the following lines:
# ExecStart=
# ExecStart=/usr/bin/oc-rsyncd --no-detach --config=/etc/oc-rsyncd.conf
Restart=on-failure
RestartSec=2s
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=true
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_FOWNER CAP_CHOWN CAP_DAC_OVERRIDE
AmbientCapabilities=CAP_DAC_READ_SEARCH CAP_FOWNER CAP_CHOWN CAP_DAC_OVERRIDE
RestrictNamespaces=yes

[Install]
WantedBy=multi-user.target
